# ops/webhook-server.service
[Unit]
Description=Microservice Stock Webhook Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/home/bxgh/microservice-stock
ExecStart=/usr/bin/python3 /home/bxgh/microservice-stock/ops/webhook_server.py
Restart=always
RestartSec=10
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
